# product/a2_ultra/addons/orin/migration/service_update/agibot_service/agibot_check_and_install.service
[Unit]
Description=Agibot Check and Install Service

[Service]
LimitCORE=infinity
Type=oneshot
ExecStart=/usr/local/bin/agibot_check_and_install.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
