# product/a2_ultra/addons/x86_64/migration/service_update/fstrim/fstrim-daily-boot.timer
[Unit]
Description=Run fstrim once per day (first boot only)

[Timer]
OnBootSec=1min
OnUnitInactiveSec=1day
AccuracySec=1h
Persistent=false

[Install]
WantedBy=timers.target
